# model-tag-monitor/deploy/ai-incubation/templates/model-tag-monitor-refresh.service
[Unit]
Description=model-tag-monitor daily dashboard refresh
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/opt/soft/model-tag-monitor/conf/model-tag-monitor.env
WorkingDirectory=/opt/soft/model-tag-monitor/app/model-tag-monitor
ExecStart=/bin/bash -lc './scripts/refresh-dashboard-daily.sh'
StandardOutput=append:/opt/log/model-tag-monitor/worker/refresh.log
StandardError=append:/opt/log/model-tag-monitor/worker/refresh.err.log
User=root
Group=root

[Install]
WantedBy=multi-user.target
